# com.user.wamacreminder.popup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.wamacreminder.popup</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/jirisaaan/Documents/Trial Project Task Reminder/show_popup.js</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/jirisaaan/Documents/Trial Project Task Reminder</string>

  <!-- Run every 5 minutes; script checks internally if past 8 AM and not shown today -->
  <key>StartInterval</key>
  <integer>300</integer>

  <key>StandardOutPath</key>
  <string>/Users/jirisaaan/Documents/Trial Project Task Reminder/popup.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/jirisaaan/Documents/Trial Project Task Reminder/popup.error.log</string>
</dict>
</plist>
